# roles/influxdb/templates/influxdb.service.j2
[Unit]
After=network-online.target

[Service]
User=root
ExecStart=

[Install]
WantedBy=multi-user.target